# gunicorn.conf
[program:gunicorn]
command=gunicorn --workers=3 app:app -b localhost:5000
directory=/deploy/app
